# debian/postinst
#!/bin/bash

set -e

# enable the service for the next boot
systemctl enable cloudshell-lcd.service
